# templates/sing-box.service
[Unit]
Description=sing-box service
Documentation=https://sing-box.sagernet.org
After=network.target nss-lookup.target network-online.target

[Service]
User=root
ExecStart=/usr/bin/sing-box run -c /etc/sing-box/config.json
ExecStartPost=/bin/bash -c 'for i in $(seq 1 10); do ip link show singbox-tun >/dev/null 2>&1 && break; sleep 1; done'
ExecStartPost=-/usr/sbin/sysctl -w net.ipv6.conf.singbox-tun.disable_ipv6=1
ExecStartPost=-/usr/bin/resolvectl dns singbox-tun ""
ExecStartPost=-/usr/bin/resolvectl domain singbox-tun ""
ExecStartPost=-/bin/bash -c '/usr/sbin/iptables -C FORWARD -o singbox-tun -j ACCEPT 2>/dev/null || /usr/sbin/iptables -I FORWARD -o singbox-tun -j ACCEPT'
ExecStartPost=-/bin/bash -c '/usr/sbin/iptables -C FORWARD -i singbox-tun -j ACCEPT 2>/dev/null || /usr/sbin/iptables -I FORWARD -i singbox-tun -j ACCEPT'
ExecStopPost=-/bin/bash -c '/usr/sbin/iptables -C FORWARD -o singbox-tun -j ACCEPT 2>/dev/null && /usr/sbin/iptables -D FORWARD -o singbox-tun -j ACCEPT || true'
ExecStopPost=-/bin/bash -c '/usr/sbin/iptables -C FORWARD -i singbox-tun -j ACCEPT 2>/dev/null && /usr/sbin/iptables -D FORWARD -i singbox-tun -j ACCEPT || true'

Restart=on-failure
RestartSec=10s
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
